# socket/crontab
0,30 * * * * curl -s http://server:8000/push/send
 